# packaging/deb/prerm
#!/bin/sh
# ClawOps Mobile Gateway — prerm
set -e

case "$1" in
  remove|deconfigure)
    systemctl stop clawops-asterisk.service    || true
    systemctl disable clawops-asterisk.service || true
    systemctl stop clawops-agent.service       || true
    systemctl disable clawops-agent.service    || true
    ;;
esac
exit 0
